# hardware/tachyon-hub-py/systemd/beekeeper-panel-py.service
[Unit]
Description=BeeKeeper M1 status LED daemon
Documentation=https://github.com/Ropehawn/beekeeper/tree/main/hardware/tachyon-hub-py
# Drives the M1 carrier's RGB status LED over I²C via the ADP8866 chip at
# 0x27 on /dev/i2c-1. Indicates Tachyon hub health (healthy / hub-down /
# ble-down / stale-upload). Does NOT touch the M1 user button — that
# button is hardwired in parallel with the SoM power-key path and has no
# userspace event semantics on this hardware.
#
# See ../ADP8866_NOTES.md for register-level gotchas (especially the
# inverted PWR_SEL polarity on this carrier).
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=600
StartLimitBurst=5

[Service]
Type=simple
User=particle
Group=particle
# i2c group lets us drive the ADP8866 over /dev/i2c-1 without root.
SupplementaryGroups=i2c
WorkingDirectory=/home/particle/beekeeper-ai
ExecStart=/home/particle/beekeeper-ai/bin/python m1_panel.py
Restart=always
RestartSec=10
Environment=PYTHONUNBUFFERED=1
StandardOutput=journal
StandardError=journal
TimeoutStopSec=15

[Install]
WantedBy=multi-user.target
