# hardware/tachyon-hub-py/systemd/beekeeper-panel-py.service
[Unit]
Description=BeeKeeper M1 panel daemon (button + ADP8866 RGB LED status)
Documentation=https://github.com/Ropehawn/beekeeper/tree/main/hardware/tachyon-hub-py
# Listens to the M1 enclosure user button on /dev/input/event2 (pmic_pwrkey
# — the M1 button is wired in parallel with the Tachyon power button).
# A logind drop-in disables HandlePowerKey so presses don't shut us down.
#
# Drives the three M1 RGB LEDs over I²C via the ADP8866 chip at 0x27 on
# /dev/i2c-1. See ../ADP8866_NOTES.md for register-level gotchas
# (especially the inverted PWR_SEL polarity on this carrier).
After=network-online.target systemd-logind.service
Wants=network-online.target
StartLimitIntervalSec=600
StartLimitBurst=5

[Service]
Type=simple
User=particle
Group=particle
# `input` group: read /dev/input/event* for the button.
# `i2c`   group: drive the ADP8866 over /dev/i2c-1.
SupplementaryGroups=input i2c
WorkingDirectory=/home/particle/beekeeper-ai
ExecStart=/home/particle/beekeeper-ai/bin/python m1_panel.py
Restart=always
RestartSec=10
Environment=PYTHONUNBUFFERED=1
StandardOutput=journal
StandardError=journal
TimeoutStopSec=15

[Install]
WantedBy=multi-user.target
